# config/jobmanager.conf
# Upstart Job to execute the Spring Boot application on startup
description "jobmanager"
start on startup
task
exec /vagrant/jobmanager/config/spring-start.sh
